# files/etc/service/php-fpm/run
#!/bin/bash
exec /usr/local/sbin/php-fpm -F
